# templates/geoipupdate-update.timer.erb
[Unit]
Description=Timer required by geoipupdate-update.service.erb
Wants=geoipupdate-update.service

[Timer]
Unit=geoipupdate-update.service
OnCalendar=<%= @timer_oncalendar %>
Persistent=true

[Install]
WantedBy=timers.target
