# templates/geoipupdate-update.timer.erb
[Unit]
Description=Timer required by geoipupdate-update.service
Wants=geoipupdate-update.service

[Timer]
Unit=geoipupdate-update.service
OnCalendar=<%= @timer_oncalendar %>

[Install]
WantedBy=timers.target
