# templates/geoipupdate-update.service.erb
# Service that updates MaxMind's 'geoipupdate' client

[Unit]
Description='geoipupdate' update service
After=network.target

[Service]
Type=simple
ExecStart=<%= @package_name %><% if @verbose %> -v<% end %> -f <%= @conf_dir %>/GeoIP.conf
User=<%= @user %>

[Install]
WantedBy=network-online.target
