# templates/geoipupdate-update.service.erb
# Service that updates MaxMind's 'geoipupdate' client

[Unit]
Description='geoipupdate' update service
After=network.target

[Service]
Type=simple
ExecStart=<%= @package_name %> -f <%= @conf_dir %>/GeoIP.conf -d <%= @target_dir %>

[Install]
WantedBy=network-online.target
